# build/utilities.conf
[program:rclone-script]
autorestart = false
startsecs = 0
user = nobody
command = /home/nobody/rclone.sh
umask = 000
