# build/utilities.conf
[program:rclone-script]
user = nobody
command = /home/nobody/rclone.sh
umask = 000

[program:logrotate-rclone-script]
user = nobody
command = /home/nobody/logrotate.sh "/config/rclone/logs/rclone.log"
umask = 000
